# scripts/myrmex-backup.service
[Unit]
Description=Myrmex Control — бэкап myrmex.json

[Service]
Type=oneshot
ExecStart=/bin/bash /root/LabDoctorM/projects/myrmex-control/scripts/backup.sh
WorkingDirectory=/root/LabDoctorM/projects/myrmex-control
User=root